# etc/systemd/system/docker.service.d/http-proxy.conf
# https://docs.docker.com/config/daemon/systemd/

[Service]
Environment="HTTP_PROXY=http://10.77.8.70:8080/"

